# gladys/rootfs/etc/services.d/gladys/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Cleanup Gladys Assistant on stop
# ==============================================================================

bashio::log.info "Gladys Assistant stopped."
