# gladys/rootfs/etc/services.d/gladys/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Gladys Assistant
# ==============================================================================

declare server_port
declare timezone
declare gladys_data_path

# Read configuration options
server_port=$(bashio::config 'server_port')
timezone=$(bashio::config 'timezone')

# Data path for Gladys persistent storage
gladys_data_path="/data/gladys"

# Create data directory if it doesn't exist
mkdir -p "${gladys_data_path}"

bashio::log.info "Starting Gladys Assistant..."
bashio::log.info "Server port: ${server_port}"
bashio::log.info "Timezone: ${timezone}"
bashio::log.info "Data path: ${gladys_data_path}"

# Check if a Gladys container already exists
if docker inspect gladys > /dev/null 2>&1; then
    bashio::log.info "Removing existing Gladys container..."
    docker rm -f gladys > /dev/null 2>&1 || true
fi

# Pull the latest Gladys image
bashio::log.info "Pulling latest Gladys Assistant image..."
docker pull gladysassistant/gladys:v4

# Start Gladys container using the official Docker run command
bashio::log.info "Starting Gladys Assistant container..."
docker run \
    --rm \
    --log-driver json-file \
    --log-opt max-size=10m \
    --cgroupns=host \
    --privileged \
    --network=host \
    --name gladys \
    -e NODE_ENV=production \
    -e SERVER_PORT="${server_port}" \
    -e TZ="${timezone}" \
    -e SQLITE_FILE_PATH=/var/lib/gladysassistant/gladys-production.db \
    -v /var/run/docker.sock:/var/run/docker.sock \
    -v "${gladys_data_path}":/var/lib/gladysassistant \
    -v /dev:/dev \
    -v /run/udev:/run/udev:ro \
    gladysassistant/gladys:v4
