# gladys/rootfs/etc/services.d/gladys/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Gladys Assistant
# ==============================================================================

declare server_port
declare timezone

# Read configuration options
server_port=$(bashio::config 'server_port')
timezone=$(bashio::config 'timezone')

bashio::log.info "Starting Gladys Assistant..."
bashio::log.info "Server port: ${server_port}"
bashio::log.info "Timezone: ${timezone}"
bashio::log.info "Data path: /data"

# Symlink Docker socket so Gladys can find it at the expected path
if [ -S /run/docker.sock ] && [ ! -e /var/run/docker.sock ]; then
    ln -s /run/docker.sock /var/run/docker.sock
fi

# Set environment variables for Gladys
export NODE_ENV=production
export SERVER_PORT="${server_port}"
export TZ="${timezone}"
export SQLITE_FILE_PATH=/data/gladys-production.db

# Start Gladys (source code is at /src/server in the Gladys image)
cd /src/server || bashio::exit.nok "Gladys source directory not found"
exec node index.js
